# bin/webchange.service
[Unit]
Description=Webchange

[Service]
User=www-data

# The configuration file application.properties should be here:
#change this to your workspace
WorkingDirectory=/srv/www/webchange

#path to executable.
#executable is a bash script which calls jar file
ExecStart=/srv/www/webchange/run

TimeoutStopSec=10
Restart=on-failure
RestartSec=5

#the code used for local update
RestartForceExitStatus=42

[Install]
WantedBy=multi-user.target
